# rust-iot.service
[Unit]
Description=Iot Certificate Authority and Iot Manager

[Service]
User=pki
Group=pki
ExecStart=rust-iot
LoadCredentialEncrypted=foobar:/etc/rust-iot/credentials

[Install]
WantedBy=multi-user.target
