# rust-iot.service
[Unit]
Description=Iot Certificate Authority and Iot Manager

[Service]
ExecStart=/usr/bin/rust-iot

[Install]
WantedBy=multi-user.target
